# deploy/hearth-fluxer.service.template
[Unit]
Description=hearth-fluxer (Fluxer presence bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__NODE_BIN__ dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
